# munin-memory.service
# NOTE: Replace <user> and <install-dir> with your system username and
# install location before installing this service file.
[Unit]
Description=Munin Memory MCP Server
After=network.target

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/<install-dir>
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=MUNIN_TRANSPORT=http
Environment=MUNIN_HTTP_PORT=3030
Environment=MUNIN_HTTP_HOST=127.0.0.1
EnvironmentFile=/home/<user>/<install-dir>/.env

# Sandboxing
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/<user>/.munin-memory
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
