# munin-offsite.service
[Unit]
Description=Munin Memory offsite backup (encrypted push to cloud)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=<user>
# Runs from ~/munin-ops (a dedicated ops dir, NOT a git checkout) so cron never
# executes code out of a mutable/dev checkout. Populate it with scripts/install-ops.sh.
WorkingDirectory=<ops-dir>
ExecStart=<ops-dir>/scripts/offsite-snapshot.sh
# Optional — reuses HEIMDALL_HUB_URL + HEIMDALL_FLEET_TOKEN (shared) for the status
# panel, plus any MIMIR_OFFSITE_* / MUNIN_OFFSITE_* overrides. Remote credentials
# and the crypt key live in rclone's chmod-600 config, not this env file or repo.
EnvironmentFile=-<ops-dir>/.env
Nice=10
IOSchedulingClass=idle

# Timeouts (VACUUM INTO snapshot + encrypted upload). Generous so a slow/throttled
# uplink can't get SIGKILLed mid-run (a kill wouldn't get to push a fail panel). The
# memory DB is small, so this is far above any real run; the heartbeat stamp
# (MIMIR_OFFSITE_STAMP) is the staleness backstop if a run is ever reaped.
TimeoutStartSec=7200

# Sandboxing
ProtectSystem=strict
# rclone may rewrite rclone.conf when it refreshes a provider OAuth token, so its
# config dir must be writable. The live DB, VACUUM INTO snapshot staging, log + stamp
# all live under ~/.munin-memory (the snapshot reads the DB; sqlite may touch WAL/shm).
ReadWritePaths=<home-dir>/.munin-memory <home-dir>/.config/rclone
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
